# config/launchd/com.newplatform.backend.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.newplatform.backend</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/sanjeev/Sites/IngestIQ/backend/.venv/bin/python</string>
    <string>-m</string>
    <string>uvicorn</string>
    <string>app.main:app</string>
    <string>--host</string>
    <string>127.0.0.1</string>
    <string>--port</string>
    <string>8889</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/sanjeev/Sites/IngestIQ/backend</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PADDLE_PDX_DISABLE_MODEL_SOURCE_CHECK</key>
    <string>True</string>
    <key>OMP_NUM_THREADS</key>
    <string>4</string>
    <key>MKL_NUM_THREADS</key>
    <string>4</string>
  </dict>
  <key>StandardOutPath</key>
  <string>/tmp/newplatform-backend.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/newplatform-backend.err.log</string>
</dict>
</plist>
